# docker/local/consumers.ini
[program:messenger-consume-async]
command=/var/www/html/bin/console messenger:consume async --time-limit=3600 -n
process_name=%(program_name)s
numprocs=1
directory=/var/www/html
umask=022
priority=999
autostart=true
autorestart=unexpected
startsecs=600
startretries=3
exitcodes=0,2
stopsignal=QUIT
stopwaitsecs=10
killasgroup=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:messenger-consume-failed]
command=/var/www/html/bin/console messenger:consume failed --time-limit=3600 -n
process_name=%(program_name)s
numprocs=1
directory=/var/www/html
umask=022
priority=999
autostart=true
autorestart=unexpected
startsecs=600
startretries=3
exitcodes=0,2
stopsignal=QUIT
stopwaitsecs=10
killasgroup=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
